# docker-production/laravel-scheduler.cron
# docker-production/laravel-scheduler.cron
* * * * * cd /var/www && php artisan schedule:run >> /dev/null 2>&1
